# panel/gx10-panel.service
[Unit]
Description=GX10 Cluster Control Panel
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
# __USER__ and __DIR__ are substituted by install.sh
User=__USER__
WorkingDirectory=__DIR__
Environment=GX10_KIT_DIR=__KIT_DIR__
Environment=GX10_PANEL_HOST=0.0.0.0
Environment=GX10_PANEL_PORT=8080
ExecStart=__DIR__/.venv/bin/python __DIR__/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
